# hyperliquid-terminal/hyperliquid-liquidations.service
[Unit]
Description=Hyperliquid Liquidations Monitor
After=network.target

[Service]
Type=simple
User=your_username
WorkingDirectory=/path/to/hyperliquid-terminal
Environment="PATH=/usr/bin:/usr/local/bin"

# Run with $50k threshold
ExecStart=/usr/bin/python3 /path/to/hyperliquid-terminal/liquidations_monitor_advanced.py 50000

# Or use this for custom threshold:
# ExecStart=/usr/bin/python3 /path/to/hyperliquid-terminal/liquidations_monitor_advanced.py 25000

# Restart on failure
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/hyperliquid-liquidations.log
StandardError=append:/var/log/hyperliquid-liquidations-error.log

[Install]
WantedBy=multi-user.target
