# hyperliquid-terminal/hyperliquid-real-liquidations.service
[Unit]
Description=Hyperliquid Real Liquidations Monitor
Documentation=https://github.com/hyperliquid/terminal
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/engine/project/hyperliquid-terminal
ExecStart=/usr/bin/python3 /home/engine/project/hyperliquid-terminal/real_liquidations_monitor.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hyperliquid-real-liquidation-monitor

[Install]
WantedBy=multi-user.target